# images/supervisor/supervisor.ini
[supervisord]
nodaemon=true
logfile=/var/www/storage/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisord.pid

[program:php-fpm]
directory=/var/www
command=php-fpm
priority=8
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:database-seed]
directory=/var/www
command=php artisan migrate:fresh --seed
autostart=true
priority=10

[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www
priority=5
numprocs=8
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/www/storage/logs/queue.log
stderr_logfile=/var/www/storage/logs/queue.log.error
stopwaitsecs=3600

[program:laravel-shedule]
process_name=%(program_name)s_%(process_num)02d
command=php artisan schedule:work
directory=/var/www
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
stdout_logfile=/var/www/storage/logs/worker.log
stderr_logfile=/var/www/storage/logs/worker.log.error
stopwaitsecs=3600